# php7/supervisor.d/crond.ini
[program:crond]
command=crond -f -l 8
autostart=true
autorestart=true
user=www-data
stderr_logfile=/var/log/cron.err.log
stdout_logfile=/var/log/cron.out.log
